# etc/profile.d/profile.sh
if [ -f "$HOME/.config/ash/profile" ]; then
	. "$HOME/.config/ash/profile"
fi
